# systemd/mastodon-backup.timer
[Unit]
Description=Run foo weekly

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
